# ansible/roles/fileserver-blackhole/files/zfs-snapshot-weekly.service
[Unit]
Description=Weekly ZFS snapshot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zfs-weekly-snapshot.sh
